# crontab
@reboot python3 ~/Ent/src/main.py
0 0 * * * python3 ~/Ent/src/main.py
